# activity-manager-mcp.service
[Unit]
Description=Activity Manager MCP Server (Streamable HTTP)
After=network.target

[Service]
Type=simple
User=activity-manager
Group=activity-manager
WorkingDirectory=/opt/activity-manager
ExecStart=/opt/activity-manager/venv/bin/python /opt/activity-manager/mcp_server/server.py --transport http
EnvironmentFile=/opt/activity-manager/.env
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/activity-manager
PrivateTmp=true

[Install]
WantedBy=multi-user.target
